# deploy/cloudflared-quick.service
[Unit]
Description=Cloudflare Quick Tunnel for donut
After=network-online.target donut.service
Requires=donut.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
ExecStart=/usr/local/bin/cloudflared tunnel --url http://127.0.0.1:8765 --no-autoupdate
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target