# deploy/donut.service
[Unit]
Description=どーなつ stream-translator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/donut/server
EnvironmentFile=-/home/ubuntu/donut/server/.env
Environment=RELOAD=false
Environment=HOST=0.0.0.0
Environment=PORT=8765
# 0.0.0.0: NPM(Docker) から 172.17.0.1 経由で到達するため。8765 は ufw で外に開けないこと
ExecStart=/home/ubuntu/donut/server/.venv/bin/uvicorn main:app --host 0.0.0.0 --port 8765
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target